# files/etc/systemd/system/cmd-pprof.timer
[Unit]
Description=Run cmd-pprof.service every 20 minutes

[Timer]
OnCalendar=*:0/20
Persistent=true

[Install]
WantedBy=multi-user.target
