# docker_cfg_files/docker-cleanup.timer
[Unit]
Description=Run docker-clenaup.service every 10 minutes

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=multi-user.target
